# amarket/offline/64-qt-js-webview-offline-installer--u18.04-qt5.2.1-r00/app/data/systemctl/qdjs.service
[Unit]
Description=Run mount RAM folder

[Service]
Type=oneshot
ExecStart=/opt/qt-desktop-js/mountfs.sh

[Install]
WantedBy=multi-user.target
